# supervisor/leaderless.supervisord.ini
[supervisord]
logfile=./supervisord.log    ; (main log file;default $CWD/supervisord.log)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=./supervisord.pid    ; (supervisord pidfile;default supervisord.pid)
nodaemon=true                ; (start in foreground if true;default false)

[program:sync-1]
command=./prefix-log.sh node ../leaderless/sync.js
environment=WAMP_PORT=9021
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=1

[program:sync-2]
command=./prefix-log.sh node ../leaderless/sync.js
environment=WAMP_PORT=9022
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=1

[program:sync-3]
command=./prefix-log.sh node ../leaderless/sync.js
environment=WAMP_PORT=9023
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=1

[program:entry-1]
command=./prefix-log.sh node ../leaderless/entry.js
environment=FOX_PORT=1735,WAMP_PORT=9031,MQTT_PORT=1881,NODE_ID=E1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=2

[program:entry-2]
command=./prefix-log.sh node ../leaderless/entry.js
environment=FOX_PORT=1736,WAMP_PORT=9032,MQTT_PORT=1882,NODE_ID=E2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=2

[program:entry-3]
command=./prefix-log.sh node ../leaderless/entry.js
environment=FOX_PORT=1737,WAMP_PORT=9033,MQTT_PORT=1883,NODE_ID=E3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=2

[program:ndb-1]
command=./prefix-log.sh node ../leaderless/ndb.js
environment=DB_FILE=./x1.db
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=3

[program:ndb-2]
command=./prefix-log.sh node ../leaderless/ndb.js
environment=DB_FILE=./x2.db
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=3

[program:ndb-3]
command=./prefix-log.sh node ../leaderless/ndb.js
environment=DB_FILE=./x3.db
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=3
